# infra/thinkaloud.service
[Unit]
Description=ThinkAloud.ai Backend
After=network.target postgresql.service redis.service

[Service]
User=root
WorkingDirectory=/root/thinkaloudai-backend/backend
ExecStart=/root/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 4
Restart=always
RestartSec=5
EnvironmentFile=/root/thinkaloudai-backend/backend/.env

[Install]
WantedBy=multi-user.target
